# sysklogd/sysklogd-1.4.1/files/DEBIAN/prerm
#!/bin/bash
# Place in this file any commands you would like to be run
# before the package has been removed

# Remove rc.d links
update-rc.d -f sysklogd remove
